# Agent/fileguard.service
[Unit]
Description=FileGuard IM Agent - Integrity Monitor
Documentation=https://github.com/KN-IMS/KN-IMS
After=network.target
Wants=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/agent -f -c /etc/im_monitor/im.conf
Restart=on-failure
RestartSec=5s
WatchdogSec=30s
TimeoutStopSec=10s
User=root
KillMode=control-group

# 보안 강화
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/run /var/log

[Install]
WantedBy=multi-user.target
