# Agent/configs/im_monitor.service
[Unit]
Description=IM Monitor — Integrity Monitoring Agent
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/agent -f -c /etc/im_monitor/im.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
WatchdogSec=90s
NotifyAccess=main

[Install]
WantedBy=multi-user.target
